# dotfiles/systemd/user/hyprdynamicmonitors.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/hyprdynamicmonitors run --enable-lid-events
